# wind_teller.service
[Unit]
Description=Wind Teller
Wants=
After=network.target

[Service]
Type=notify

User=yggdra
Group=yggdra

ExecStart=/usr/local/bin/wind_teller

Restart=always

TimeoutStartSec=15
WatchdogSec=10

[Install]
WantedBy=multi-user.target
